# groundstation/autostart/foxtelem.desktop
[Desktop Entry]
Type=Application
Name=FoxTelem
Terminal=true
Exec=/home/pi/CubeSatSim/groundstation/fox-startup.sh
