# groundstation/autostart/foxtelem.desktop
[Desktop Entry]
Type=Application
Name=FoxTelem
Exec=/home/pi/CubeSatSim/groundstation/fox-startup.sh
